# .context/systemd/termlink-framework-agent.service
# TermLink Framework Agent — persistent named session for cross-host framework operations
# Source of truth: /opt/termlink/.context/systemd/termlink-framework-agent.service (git-tracked)
# Installed to:    /etc/systemd/system/termlink-framework-agent.service (copy)
# Install with:    sudo .context/systemd/install.sh
#
# Provides a shell-backed, discoverable session that other agents on the ring20 LAN
# can target via `termlink agent ask`, `termlink remote exec`, or `termlink remote push`.
# Handles: fw pickup process, fw audit, fw task list, and general framework operations.

[Unit]
Description=TermLink Framework Agent — persistent shell session (role: framework)
After=termlink-hub.service
Wants=termlink-hub.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/opt/termlink

Environment=TERMLINK_RUNTIME_DIR=/var/lib/termlink
Environment=PROJECT_ROOT=/opt/termlink

ExecStart=/root/.cargo/bin/termlink register \
  --name framework-agent \
  --roles framework,pickup \
  --tags host=107,project=termlink \
  --cap data_plane,stream \
  --shell \
  --allowed-commands "fw,cat,ls,head,grep,tail,wc,date,python3"

Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/termlink /opt/termlink/.context /opt/termlink/.tasks /tmp

[Install]
WantedBy=multi-user.target
